# config/examples/etc/systemd/services/torrent-hdd.service
[Unit]
Description=Torrent start
Requires=docker.service
After=torrent.service

[Service]
Type=oneshot
TimeoutStartSec=0

WorkingDirectory=/home/dima/media

ExecStartPre=/usr/bin/docker-compose -f docker-compose-torrent-hdd.yml down
ExecStart=/usr/bin/docker-compose -f docker-compose-torrent-hdd.yml up

[Install]
WantedBy=default.target
