# config/examples/etc/systemd/services/torrent.service
[Unit]
Description=Torrent start
Requires=docker.service
After=torrent-pull.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/dima/media
ExecStart=/usr/bin/docker-compose -f docker-compose-torrent.yml up
ExecStop=/usr/bin/docker-compose -f docker-compose-torrent.yml down
TimeoutStartSec=0

[Install]
WantedBy=default.target
